# etc/systemd/system/gunicorn_shop.socket
[Unit]
Description=gunicorn socket

[Socket]
ListenStream=/run/gunicorn_shop.sock

[Install]
WantedBy=sockets.target
